# data/com.example.Gtk.JSApplication.service.in
[Unit]
Description=JS Application
Documentation=http://www.example.com/gtk-js-app/

[Service]
BusName=com.example.Gtk.JSApplication
ExecStart=@pkgdatadir@/com.example.Gtk.JSApplication
#Type=application
ExecStop=@bindir@/gapplication com.example.Gtk.JSApplication quit
